# board/systemd/gemy-autostart.service
[Unit]
Description=Gemy greeting robot (auto-start on power-on)
After=multi-user.target sound.target
Wants=gemy-watcher.service

[Service]
Type=simple
ExecStart=/home/root/gemy-boot.sh
# Do not restart on crash (avoids buzzer/log storm while tuning).
Restart=no
WorkingDirectory=/home/root

[Install]
WantedBy=multi-user.target
